# roles/linux-patch/templates/systemd-auto-updates.timer.j2
# DO NOT EDIT THIS FILE - RISK OF CONFIGURATION INCONSISTENCY
# DEPLOY ANY CHANGES WITH ANSIBLE FROM THE CORRESPONDING BITBUCKET REPOSITORY
#
# Timer für die automatischen Update
# Autor: Philip Knack
 
[Unit]
Description=System Updates automatisch installieren
Requires=auto-updates.service
Documentation=https://wiki.materna.de/display/ZEMSO/ECE+Automatisierungskonzept+-+Automatisiertes+Linux+patchen
 
[Timer]
# https://manpages.debian.org/testing/manpages-de/systemd.time.7.de.html
Unit=auto-updates.service
# Die Uhrzeit ist zufällig gewählt und sollte (hoffentlich) keine anderen Tasks stören.
OnCalendar=Mon *-*-* 01:25:00
 
[Install]
WantedBy=timers.target